# build/supervisord.conf
[supervisord]
nodaemon=true

[program:headscale]
command=/headscale serve

[program:flyctl]
command=/flyctl proxy 7000:7000 -a %(ENV_FLY_APP_NAME)s -t %(ENV_FLY_ACCESS_TOKEN)s

[program:frpc]
command=/frpc -c /etc/frpc/frpc.ini